# etc/profile.d/disable-beep.sh
case "$TERM" in 
	linux)
		setterm --blength 0
		;;
	*)
esac
